# etc/supervisord.conf
[supervisord]
pidfile=/tmp/supervisord.pid
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=true

[program:app-nginx]
stdout_events_enabled = true
stderr_events_enabled = true
command = nginx -c /src/etc/nginx.conf

[program:app]
stdout_events_enabled = true
stderr_events_enabled = true
stdout_logfile = /var/log/eebook/gunicorn_stdout.log
stderr_logfile = /var/log/eebook/gunicorn_stderr.log
redirect_stderr = true
command = gunicorn -c /src/etc/gunicorn_cfg.py wsgi:app
